# deploy/trip_planner_ia.service
[Unit]
Description=Trip Planner IA Phoenix
After=network.target

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/opt/trip_planner_ia/current
EnvironmentFile=/etc/trip_planner_ia/env
ExecStart=/opt/trip_planner_ia/current/bin/server
Restart=always
RestartSec=5
MemoryMax=450M
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target